# 1-python3-uwsgi-nginx-with-0/keep/supervisor/app-nginx.conf
[program:app-nginx]
command = /usr/sbin/nginx
priority = 5500
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
