# 2-django-postgresql-with-1/run/home-python3/conf-dev/supervisor/app-djangosrv.conf
[program:app-djangosrv]
command = /home/python3/conf-dev/djangosrv.sh
priority = 5000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
