# ubuntu-python3-uwsgi-nginx/keep/conf.prod/supervisor.conf
[program:app-uwsgi]
command = /home/python3/envs/app/bin/uwsgi --ini /home/python3/uwsgi.ini

[program:app-nginx]
command = /usr/sbin/nginx
