# 2-django-postgresql-with-1/keep/supervisor-dev/app-postgresql.conf
[program:app-postgresql]
command = /etc/init.d/postgresql start
priority = 4500
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
